# deploy/tb-nightly.service
[Unit]
Description=toronto-bids nightly sync, enrich, export, notify, publish
Documentation=https://github.com/CivicTechTO/toronto-bids/blob/main/docs/superpowers/specs/2026-07-19-publish-data-design.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/toronto-bids/scrapers
# The default data dir is scrapers/files/ — inside the checkout. On a server the archive must
# live where a `git pull` cannot reach it.
Environment=TB_DATA_DIR=%h/tb-data
# The leading '-' makes this optional: no webhook file -> notify.post() no-ops and the run
# still succeeds. The file holds credentials (TB_SLACK_WEBHOOK, optional GH_TOKEN) and is mode
# 0600; it is never in git. Publishing requires gh to be authenticated (a GH_TOKEN here OR a
# prior `gh auth login` on the box) — the publish step fails loudly otherwise (#146).
EnvironmentFile=-%h/.config/toronto-bids/tb.env
# The wrapper runs `tb nightly` then publish-data.sh — publish runs even after a partial sync
# (the export is still valid) and the unit fails if EITHER step failed, neither masking the
# other. See docs/superpowers/specs/2026-07-19-publish-data-design.md (#146).
ExecStart=%h/toronto-bids/deploy/tb-nightly-run.sh
# A full sync measures ~3 minutes (sync_run, 2026-07-16); the upload adds a little. 30m ceiling.
TimeoutStartSec=30m
